# docker/supervisor.conf
[supervisord]
nodaemon=true

[program:markdown_api]
directory=/node_app
command=/node_app/docker/start.sh
redirect_stderr=true
autostart=true
autorestart=true
